# inotify/cadratin.service
[Service]
ExecStart=<project>/inotify/cadratin_inotify.sh
